# roles/pub_storm/files/nimbus.ini
[program:nimbus]
command=/opt/apache-storm-0.9.4/bin/storm nimbus
user=root
autostart=false
autorestart=true
startsecs=3
