# roles/pub_storm/files/ui.ini
[program:ui]
command=/opt/apache-storm-0.9.4/bin/storm ui
user=root
autostart=false
autorestart=true
startsecs=3