# roles/pub_tinc/files/tinc.ini
[program:tinc]

command=tincd -n vpn
user=root

autostart=false

autorestart=true

startsecs=3
