# roles/pub_zookeeper/files/zookeeper.ini
[program:zookeeper]
command=/opt/zookeeper-3.4.6/bin/zkServer.sh start-foreground
user=root
autostart=false
autorestart=true
startsecs=3

