# roles/pub_redis/files/redis.ini
[program:redis]

command = redis-server

user=root

autostart=false

autorestart=true

startsecs=3
