# roles/pub_mysql/files/mysql.ini
[program:mysql]

command=mysql

user=root

autostart=false

autorestart=true

startsecs=3
