# roles/pub_hadoop/files/jobtracker.ini
[program:jobtracker]

command=/opt/hadoop-2.6.3/sbin/hadoop-daemons.sh start jobtracker
user=root

autostart=false

autorestart=true

startsecs=3
