# roles/pub_hadoop/files/namenode.ini
[program:namenode]

command=/opt/hadoop-2.6.3/sbin/hadoop-daemons.sh start namenode
user=root

autostart=false

autorestart=true

startsecs=3
