# roles/pub_hadoop/files/tasktracker.ini
[program:tasktracker]

command=/opt/hadoop-2.6.3/sbin/hadoop-daemons.sh start tasktracker
user=root

autostart=false

autorestart=true

startsecs=3
